# deploy/satellite-backend.service
[Unit]
Description=Satellite Data FastAPI backend (geotiles.tech)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Change User if the app is not owned by 'ubuntu'
User=ubuntu
WorkingDirectory=/opt/satellite-data/backend
# pydantic-settings reads /opt/satellite-data/.env automatically (path is derived
# from the source tree), and systemd also injects it here for good measure.
EnvironmentFile=/opt/satellite-data/.env
ExecStart=/opt/satellite-data/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=3
# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
